# net/netkit-rsh/xinetd.d/rlogin
service login
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	group		= tty
	server		= /usr/sbin/in.rlogind
	server_args	=
	log_type	= FILE /var/log/xinetd/rlogin
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
