# net/linux-identd/xinetd.d/linux-xinetd
service auth
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= nobody
	server		= /usr/sbin/in.identd
	server_args	= 
	per_source	= 4
	log_type	= FILE /var/log/xinetd/linux-identd
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
