# net/samba/xinetd.d/netbios-ns
service netbios-ns
{
	socket_type	= dgram
	protocol	= udp
	wait		= yes
	user		= root
	group		= root
	server_args	= -o  -s /etc/samba/smb.conf
	server		= /usr/sbin/nmbd
}
