# archive/amanda/xinetd.d/amanda
service amanda
{
	socket_type	= dgram
	protocol	= udp
	wait		= yes
	user		= backup
	group		= backup
	server		= /usr/libexec/amandad
	log_type	= FILE /var/log/xinetd/amanda
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
