# news/leafnode/xinetd.d/nntp
service nntp
{
  disable = no
  flags           = REUSE
  protocol        = tcp
  socket_type     = stream
  wait            = no
  user            = news
  only_from       = localhost
  server          = /usr/sbin/leafnode
  log_on_failure  += USERID
}


