# mail/imap/xinetd.d/pop3-stunnel
service pop3s
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= root
	#user			= nobody
	server			= /usr/sbin/stunnel
	server_args		= -l /usr/sbin/ipop3d
	# pop3s will work without stunnel if you have /etc/ssl/certs/ipop3d.pem
	# server		= /usr/sbin/ipop3d
	# server_args	=
	per_source		= 4
	log_type		= FILE /var/log/xinetd/pop3s
	log_on_success		= PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
