# mail/imap/xinetd.d/pop3
service pop3
{
	socket_type		= stream
	protocol		= tcp
	wait			= no
	user			= root
	#user			= nobody
	server			= /usr/sbin/ipop3d
	#server_args		= 
	per_source		= 4
	log_type		= FILE /var/log/xinetd/pop3
	log_on_success		= PID HOST USERID EXIT DURATION
	log_on_failure		= USERID ATTEMPT
}
