# net/netkit-ntalk/xinetd.d/talk
service ntalk
{
	socket_type	= dgram
	protocol	= udp
	wait		= no
	user		= nobody
	group		= tty
	server		= /usr/sbin/in.talkd
#	server_args	=
	per_source	= 16
	log_type	= FILE /var/log/xinetd/talk
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
