# mail/qmail/xinetd.d/qmail-smtpd
service smtp
{
    socket_type             = stream
    protocol                = tcp
    wait                    = no
    user                    = qmaild
    server                  = /var/qmail/bin/tcp-env
    server_args             = /var/qmail/bin/qmail-smtpd
    log_on_success          = PID HOST USERID EXIT DURATION
    log_on_failure          = HOST USERID ATTEMPT
}
