# utils/sysstat/cron.hourly/sa1
#!/bin/sh

# Collect and store binary data in the system activity daily data file
# Run system activity accounting tool every 10 minutes

/usr/lib/sa/sa1 600 6 -S DISK &
