# chat-irc/bitlbee/xinetd.d/bitlbee
service ircd-bitlbee
{
	type		= UNLISTED
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= nobody
	server		= /usr/sbin/bitlbee
	port		= 6667
## uncomment the following if for personal use only
# 	interface	= localhost
# 	per_source	= 1
	log_type	= FILE /var/log/xinetd/bitlbee
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
