# ftp/wu-ftpd/xinetd.d/wu-ftpd
service ftp
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/sbin/in.ftpd
	server_args	= -l -a
	per_source	= 8
	log_type	= FILE /var/log/xinetd/wu-ftpd
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
