# devel/distcc/xinetd.d/distcc
service distcc
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= daemon
	server		= /usr/bin/distccd
	server_args	= --inetd
	per_source	= 8
	log_type	= FILE /var/log/xinetd/distccd
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
