# smgl/smgl-upstart/upstart-init/native/clamav.conf
# clamav
#
# starts clamd

description	"Start ClamAV"
author		"Arjan Bouter <abouter@sourcemage.org>"

emits clamav-started clamav-stopped

start on network-start and syslog-start
stop on shutdown

exec /usr/sbin/clamd
