# smgl/smgl-upstart/upstart-init/native/slapd.conf
# slapd
#
# starts slapd

description	"Start OpenLDAP slapd"
author		"Arjan Bouter <abouter@sourcemage.org>"

start on network-up and started syslog
stop on shutdown

script
  . /etc/sysconfig/slapd
  slapd -d -h $SLAPD_URL
end script
