# smgl/smgl-upstart/upstart-init/native/xinetd.conf
# xinetd
#
# starts xinetd

description	"Find and enable LVM"
author		"Arjan Bouter <abouter@sourcemage.org>"

start on network-up and syslog-start
stop on shutdown

exec /usr/sbin/xinetd -dontfork
