# cluster/openmosixview/xinetd.d/mosstatd
service mosstatd
{
	socket_type		= stream
	wait			= no
	user			= nobody
	log_on_success		+= USERID
	log_on_failure 		+= USERID
	server			= /usr/bin/mosstatd
}



