# ftp/vsftpd/xinetd.d/vsftpd
# vsftpd is the secure FTP server.
service ftp
{
        disable                 = no
        socket_type             = stream
        wait                    = no
        user                    = root
#        server                  = /usr/local/bin/vsftpd
        server                  = /usr/sbin/vsftpd
        per_source              = 5
        instances               = 200
        log_on_success          += PID HOST DURATION
        log_on_failure          += HOST
		disabled				= yes
}
