# audio-drivers/alsa-oss/profile.d/alsa-oss
# Set this variable if you want programs to use the 
# oss-to-alsa wrapper instead of the OSS kernel driver.
# USE_AOSS_WRAPPER=y
