# audio-drivers/alsa-oss/profile.d/alsa-oss.sh
#!/bin/bash

[[ -f /etc/sysconfig/alsa-oss ]]  &&
  . /etc/sysconfig/alsa-oss       &&
  export  USE_AOSS_WRAPPER
