# utils/sysstat/cron.daily/sa2
#!/bin/sh

# Write daily report in the /var/log/sa directory

/usr/lib/sa/sa2 -A &
