# mail/courier-authlib/init.d/courier-authlib
#! /bin/sh
#
# chkconfig: 2345 15 85
# description: Authlib - Courier authentication library.
#
#

RUNLEVEL=3
NEEDS="+network +remotefs"
.  /etc/init.d/smgl_functions

prefix="/usr"
exec_prefix="/usr"
sysconfdir="/etc"
sbindir="/usr/sbin"

if test ! -f /etc/courier/authlib/authdaemonrc
then
	echo "/etc/courier/authlib/authdaemonrc does not exist, forgot make install-configure?"
	exit 1
fi

case "$1" in
start)
        cd /
        # Start daemons.
        touch /var/lock/subsys/courier-authlib

        echo -n "Starting Courier authentication services:"

	$sbindir/authdaemond start

	echo " authdaemond"
        ;;
stop)
        echo -n "Stopping Courier authentication services:"

	$sbindir/authdaemond stop

	echo " authdaemond"
        ;;
restart)
	$0 stop
	$0 start
        ;;
reload)
        echo -n "Starting Courier authentication services:"

	$sbindir/authdaemond restart

	echo " authdaemond"
        ;;
esac
exit 0
