# net/sslh/xinetd.d/sslh
service sslh
{
        socket_type     = stream
        protocol        = tcp
        wait            = no
        user            = nobody
        server          = /usr/sbin/sslh-fork
        server_args     = --inetd
        log_type        = FILE /var/log/xinetd/sslh
        log_on_success  = PID HOST USERID EXIT DURATION
        log_on_failure  = USERID ATTEMPT
}
