# collab/subversion/xinetd.d/svnserve
service svnserve
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= svn
	server		= /usr/bin/svnserve
	server_args	= -i
	per_source	= 8
	log_type	= FILE /var/log/xinetd/svnserve
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
