# utils/fam/xinetd.d/fam
service sgi_fam
{
        type         =  RPC UNLISTED
        socket_type  =  stream
        user         =  root
        group        =  nogroup
        server       =  /usr/sbin/famd
        wait         =  yes
        protocol     =  tcp
        rpc_version  =  2
        rpc_number   =  391002
}
