# mail/teapop/xinetd.d/teapop
service pop3
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/libexec/teapop
#	server_args	=
	per_source	= 4
	log_type	= FILE /var/log/xinetd/teapop
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
