# mail/teapop/xinetd.d/teapop-stunnel
service pop3s
{
	socket_type	= stream
	protocol	= tcp
	wait		= no
	user		= root
	server		= /usr/sbin/stunnel
	server_args	= -l /usr/libexec/teapop
	per_source	= 4
	log_type	= FILE /var/log/xinetd/teapop-stunnel
	log_on_success	= PID HOST USERID EXIT DURATION
	log_on_failure	= USERID ATTEMPT
}
