# utils/acpid2/systemd/acpid.service
[Unit]
Description=ACPI daemon

[Service]
PIDFile=/run/acpid.pid
ExecStart=/usr/sbin/acpid -l -f

[Install]
WantedBy=multi-user.target
